# deploy/job-automate.service
[Unit]
Description=Job poller (GitHub repos + jobhive ATS scrape -> Supabase + Discord)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=rpi
Group=rpi
WorkingDirectory=/home/rpi/job-automate/engine
EnvironmentFile=/home/rpi/job-automate/.env
ExecStart=/home/rpi/job-automate/.venv/bin/python /home/rpi/job-automate/engine/run.py
# One run (per-company jitter + ATS scrape + retries) can exceed the 90s default.
TimeoutStartSec=900

# No [Install]: this is triggered by job-automate.timer (and runnable on demand
# with `systemctl start job-automate.service`). run.py stays runnable manually too.
