# deploy/job-automate-nuworks.timer
[Unit]
Description=Run the NUworks poll hourly at a randomized 45-60 min offset

[Timer]
# Base fire time: :45 past every hour, plus RandomizedDelaySec=[0, 900s] on top,
# so each run lands randomly between :45 and :60 past the hour -- i.e. a random
# 45-60 min delay. Low-volume (3 API pages) plus unpredictable timing keeps the
# polling well clear of anything that looks like automated abuse.
OnCalendar=*-*-* *:45:00
RandomizedDelaySec=900
# Catch up a run missed during reboot/downtime.
Persistent=true
Unit=job-automate-nuworks.service

[Install]
WantedBy=timers.target
