# deploy/job-automate-jobright.service
[Unit]
Description=Jobright + JobSpy scrapers (-> jobright_jobs / jobspy_jobs + Discord)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=rpi
Group=rpi
WorkingDirectory=/home/rpi/job-automate/engine
EnvironmentFile=/home/rpi/job-automate/.env
# Two scrapers share this one oneshot (run sequentially), so they share the
# jobright timer rather than each getting its own. jobright is fast; jobspy
# scrapes 4 boards and is slower, hence the larger timeout below.
ExecStart=/home/rpi/job-automate/.venv/bin/python /home/rpi/job-automate/engine/jobright.py
ExecStart=/home/rpi/job-automate/.venv/bin/python /home/rpi/job-automate/engine/jobspy_scraper.py
TimeoutStartSec=900

# No [Install]: triggered by job-automate-jobright.timer (and runnable on demand
# with `systemctl start job-automate-jobright.service`). jobright.py stays
# runnable manually too. Intentionally SEPARATE from job-automate.service --
# different table, different pipeline.
